# configs/etc/s6/sshd/run
#!/bin/bash

# setup cross compilation toolchain
/opt/setuptoolchain.sh

# setup users
/opt/setupusers.sh

# start ssh
/usr/sbin/sshd -D
