# configs/etc/s6/sshd/run
#!/bin/bash

# setup users, git and development environment
/opt/setupusers.sh
/opt/setupgit.sh
/opt/setupenv.sh

# setup cross compilation toolchain
/opt/setuptoolchain.sh

# start ssh
exec /usr/sbin/sshd -D
