# configs/etc/s6/sshd/run
#!/bin/bash

# setup cross compilation toolchain
/opt/setuptoolchain.sh

# setup users and git
/opt/setupusers.sh
/opt/setupgit.sh

# start ssh
/usr/sbin/sshd -D
